# usr/lib/systemd/system/greetd.service.d/00-nobleed.conf
[Service]
ExecStartPre=kill -SIGRTMIN+21 1
ExecStopPost=kill -SIGRTMIN+20 1

